# airootfs/etc/systemd/system/multi-user.target.wants/opsec-live-setup.service
[Unit]
Description=OpsecOS Live Setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/opsec-live-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
