# airootfs/etc/systemd/system/ghost-live-setup.service
[Unit]
Description=GhostOS Live System Setup
After=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ghost-live-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
